# autostart/lxqt-desktop.desktop.in
[Desktop Entry]
Type=Application
TryExec=lwqt-pcmanfm
Exec=/bin/sh -c "sleep 1; exec lwqt-pcmanfm --desktop --profile=lxqt"
OnlyShowIn=LWQt;
X-LXQt-Module=true

#TRANSLATIONS_DIR=translations
